# aman-pharma.service
[Unit]
Description=Aman Pharma Medical Store Server
After=network.target

[Service]
Type=simple
User=amex
WorkingDirectory=/home/amex/Documents/aman-pharma/backend
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
